# config/.config/nvim/autostart/hiddify.desktop
[Desktop Entry]
Type=Application
Name=hiddify
Comment=hiddify startup script
Exec=/opt/hiddify-next/hiddify
StartupNotify=false
Terminal=false
